# tools/whatsapp-bridge-kit/launchagents/com.patricktest.whatsapp-bridge-z2.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.patricktest.whatsapp-bridge-z2</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/patricktest/whatsapp-bridges/z2-41792168200/whatsapp-bridge/whatsapp-bridge</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/patricktest/whatsapp-bridges/z2-41792168200/whatsapp-bridge</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>DISABLE_AUTH_CHECK</key>
        <string>true</string>
        <key>API_PORT</key>
        <string>8081</string>
        <key>BRIDGE_ID</key>
        <string>z2</string>
        <key>BRIDGE_PHONE</key>
        <string>41792168200</string>
        <key>BRIDGE_ZONE</key>
        <string>formation-non-pro</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/patricktest/whatsapp-bridges/z2-41792168200/whatsapp-bridge/bridge.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/patricktest/whatsapp-bridges/z2-41792168200/whatsapp-bridge/bridge.err.log</string>
</dict>
</plist>
