# whatsapp-bridges-setup/launchd/energy.vlbh.wa-healthcheck.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>energy.vlbh.wa-healthcheck</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/patricktest/whatsapp-bridges/shared/scripts/healthcheck.sh</string>
    </array>

    <key>StartInterval</key>
    <integer>900</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/patricktest/Library/Logs/whatsapp-mcp/healthcheck.out</string>

    <key>StandardErrorPath</key>
    <string>/Users/patricktest/Library/Logs/whatsapp-mcp/healthcheck.err</string>
</dict>
</plist>
